# services/picom.service
[Unit]
Description=Compositing manager
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=/usr/bin/picom --backend glx --glx-no-stencil --vsync --shadow-exclude '! name~=""'
Environment=DISPLAY=:0

[Install]
WantedBy=graphical-session.target

