# services/wake.target
[Unit]
Description=Unlocked and whatnot
BindsTo=display-manager.service
