# services/wake.target
[Unit]
Description=Unlocked and whatnot
BindsTo=xsession.target
