# services/xsession.target
[Unit]
Description=Xsession up and running
BindsTo=graphical-session.target
